# douyin-tiktok-api.service
[Unit]
Description=Douyin TikTok Download API Service
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/douyin-tiktok-api
Environment="PATH=/opt/douyin-tiktok-api/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/douyin-tiktok-api/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 9000 --workers 4 --log-level info
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=5

# 安全设置
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/douyin-tiktok-api

# 资源限制
LimitNOFILE=65535
LimitNPROC=32768

[Install]
WantedBy=multi-user.target